# systemd/nutrientcontent.service
[Unit]
Description=nutrientcontent (NEVO microservice)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=exec
User=erik
Group=erik
WorkingDirectory=/home/erik/microservices/nutrientcontent
EnvironmentFile=/home/erik/microservices/nutrientcontent/.env
ExecStart=/home/erik/microservices/nutrientcontent/.venv/bin/uvicorn \
    src.app:app \
    --host 127.0.0.1 \
    --port 5555 \
    --no-access-log \
    --log-level info
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/erik/microservices/nutrientcontent
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true

[Install]
WantedBy=multi-user.target
